# splread.socket
[Unit]
Description=splread Output Socket

[Socket]
RemoveOnStop=true
SocketMode=0600
ListenFIFO=/home/herbivore/var/splread.fifo
